# supervisor/file_server.conf
[program:file-server]
directory=/opt/static-file-server/src
command=python3 file_server.py
